# AI_Log_Analysis-Project-1/bedrock-log-analyzer-ui/crontab
*/5 * * * * cd /app && /usr/local/bin/python auto_analyzer.py >> /var/log/auto_analyzer.log 2>&1
